# RobertLog/docker/uwsgi.conf
# simple uWSGI script
#/etc/init/uwsgi.conf

env UWSGI=/usr/local/bin/uwsgi
env LOGTO=/var/log/uwsgi/emperor.log

description "uwsgi tiny instance"
start on runlevel [2345]
stop on runlevel [06]
respawn

exec $UWSGI  /home/PlayGround/RobertLog/docker/radish_uwsgi.ini

#exec $UWSGI --master --emperor /etc/uwsgi/vassals --die-on-term --uid www-data --gid www-data --logto $LOGTO
